# example.cron
# Monitor SHIBA's price
* * * * * /scripts/crypto/shiba/check_shiba_price.sh >> /var/log/crypto/check_shiba_price.log
